# packages/wmbusmeters/wmbusmeters.service
[Unit]
Description=wmbusmeters service
After=network.target

[Service]
Type=forking
User=wmbusmeters
ExecStart=wmbusmetersd /run/wmbusmeters/wmbusmeters.pid
PIDFile=/run/wmbusmeters/wmbusmeters.pid

[Install]
WantedBy=multi-user.target
