# packages/blackarch-recon/ip2clue/ip2clued.service
[Unit]
Description=Small memory/CPU footprint daemon to lookup country (and other info) based on IP (v4 and v6)
After=network.target

[Service]
Type=forking
ExecStart=/usr/sbin/ip2clued
GuessMainPID=yes

[Install]
WantedBy=multi-user.target
