# packages/hostapd-mana/hostapd-mana.service
[Unit]
Description=Modified hostapd for Wi-Fi attacks to create a rogue access point.
After=network.target

[Service]
ExecStart=/usr/bin/hostapd /etc/hostapd-mana/hostapd.conf
ExecReload=/bin/kill -HUP $MAINPID

[Install]
WantedBy=multi-user.target
