# packages/gvmd/gvmd.service
[Unit]
Description=Greenbone Vulnerability Manager
After=ospd-openvas.service
Requires=ospd-openvas.service

[Service]
Type=simple
User=gvm
Group=gvm
PIDFile=/run/gvm/gvmd.pid
EnvironmentFile=-/etc/default/gvmd
ExecStart=/usr/bin/gvmd $OPTIONS

[Install]
WantedBy=multi-user.target
