# packages/openvas-manager/openvas-manager.service
[Unit]
Description=Open Vulnerability Assessment System Manager Daemon
Documentation=http://www.openvas.org/
Requires=openvas-scanner.service
After=openvas-scanner.service

[Service]
Type=simple
ExecStart=/usr/bin/openvasmd -f -a 127.0.0.1

[Install]
WantedBy=multi-user.target
