# packages/novnc/novnc.service
[Unit]
Description=noVNC Service
After=network.target

[Service]
Type=idle
DynamicUser=yes
Restart=on-failure
RestartSec=5s
ExecStart=/usr/bin/novnc --vnc %i

[Install]
WantedBy=multi-user.target