# packages/blackarch-defensive/snort/snort.service
[Unit]
Description=Snort IDS system listening on '%I'

[Service]
Type=simple
ExecStartPre=/usr/sbin/ip link set up dev %I
ExecStart=/usr/bin/snort -A fast -b -p -u snort -g snort -c /etc/snort/snort.conf -i %I

[Install]
Alias=multi-user.target.wants/snort@%i.service
