# packages/elasticsearch/elasticsearch-keystore.service
[Unit]
Description=Elasticsearch Keystore Generation
ConditionPathExists=|!/etc/elasticsearch/elasticsearch.keystore

[Service]
Type=oneshot
Group=elasticsearch
UMask=0007
ExecStart=/usr/share/elasticsearch/bin/elasticsearch-keystore create
RemainAfterExit=yes
