# packages/xrdp/xrdp.service
[Unit]
Description=xrdp daemon
Requires=xrdp-sesman.service
After=syslog.target network.target xrdp-sesman.service

[Service]
Type=forking
PIDFile=/var/run/xrdp.pid
EnvironmentFile=/etc/xrdp/xrdp.ini
ExecStart=/usr/bin/xrdp $XRDP_OPTIONS
ExecStop=/usr/bin/xrdp $XRDP_OPTIONS --kill

[Install]
WantedBy=multi-user.target
