# packages/xrdp/xrdp-sesman.service
[Unit]
Description=xrdp session manager
After=syslog.target network.target
StopWhenUnneeded=true

[Service]
Type=forking
PIDFile=/var/run/xrdp-sesman.pid
ExecStart=/usr/bin/xrdp-sesman
ExecStop=/usr/bin/xrdp-sesman --kill

[Install]
WantedBy=multi-user.target

