# packages/openvas-scanner/greenbone-nvt-sync.service
[Unit]
Description=Synchronize the greenbone NVT DB
Requires=network-online.target
After=network-online.target
After=redis.service

[Service]
Type=oneshot
User=gvm
Group=gvm
Restart=on-failure
RestartSec=300
ExecStart=/usr/bin/greenbone-nvt-sync
ExecStart=/usr/bin/openvas --update-vt-info

[Install]
WantedBy=multi-user.target