# packages/openvas-scanner/openvas-scanner.service
#Experimental
[Unit]
Description=Open Vulnerability Assessment System Scanner Daemon
Documentation=http://www.openvas.org/

[Service]
Type=forking
PIDFile=/run/openvas-scanner.pid
ExecStart=/usr/bin/openvassd >> /var/log/openvas-daemon.log

[Install]
WantedBy=multi-user.target
