# packages/netkit-rwho/rwhod.cron
#!/bin/sh -e
exec find /var/spool/rwho -mindepth 1 -mtime +30 -delete
