# packages/netkit-rwho/rwhod.service
[Unit]
Description=Remote WHO Daemon
After=syslog.target network.target

[Service]
Type=forking
EnvironmentFile=/etc/conf.d/rwhod
ExecStart=/usr/bin/rwhod $RWHOD_OPTS

[Install]
WantedBy=multi-user.target
