# packages/apache-tika/apache-tika.service
[Unit]
Description=Apache Tika

[Service]
Type=Simple
PIDFile=/var/run/apache-tika.pid
ExecStart=/usr/bin/tika

[Install]
WantedBy=multi-user.target
