# packages/neo4j-community/neo4j.service
[Unit]
Description=Neo4j

[Service]
User=neo4j
Type=forking
Environment=NEO4J_CONF=/etc/neo4j
Environment=NEO4J_HOME=/usr/share/java/neo4j
ExecStart=/usr/bin/neo4j start
ExecStop=/usr/bin/neo4j stop
PIDFile=/run/neo4j/neo4j.pid
LimitNOFILE=40000

[Install]
WantedBy=multi-user.target
