# archtrack/packages/sslh/sslh-select.service
[Unit]
Description=SSL/SSH multiplexer
After=network.target

[Service]
ExecStart=/usr/bin/sslh-select -f -F /etc/sslh.conf

[Install]
WantedBy=multi-user.target
