# archtrack/packages/sslh/sslh-fork.service
[Unit]
Description=SSL/SSH multiplexer
After=network.target

[Service]
ExecStart=/usr/bin/sslh-fork -f -F /etc/sslh.conf
KillMode=process

[Install]
WantedBy=multi-user.target
