# packages/ospd-openvas/ospd-openvas.service
[Unit]
Description=Openvas OSP Server
After=redis.service

[Service]
Type=simple
User=gvm
Group=gvm
EnvironmentFile=-/etc/default/ospd-openvas
ExecStartPre=-/usr/bin/sh -c 'rm -rf /run/gvm/ospd-openvas.pid /run/gvm/ospd.sock /run/gvm/gvmd.sock'
ExecStart=/usr/bin/ospd-openvas \
--pid-file=/run/gvm/ospd-openvas.pid \
--lock-file-dir=/run/gvm \
--unix-socket=/run/gvm/ospd.sock \
$OPTIONS

[Install]
WantedBy=multi-user.target
