# packages/bactl/bactl.install
post_install() {
  set -e
  systemctl enable bactl.service
  systemctl start bactl.service
}

post_upgrade() {
  post_install "$@"
}
