# packages/logkeys/logkeys.service
[Unit]
Description=Keylogger daemon

[Service]
Type=forking
GuessMainPID=yes
EnvironmentFile=/etc/conf.d/logkeysd
ExecStart=/usr/bin/logkeys --start $LOGKEYS_PARAMS
ExecStop=/usr/bin/logkeys --kill
     
[Install]
WantedBy=multi-user.target
