# packages/blackarch-proxy/cntlm/cntlm.service
[Unit]
Description=Cntlm Proxy Server
After=network.target

[Service]
Type=forking
ExecStart=/usr/bin/cntlm -P /run/cntlm.pid
PIDFile=/run/cntlm.pid
ExecStop= /usr/bin/killall -9 cntlm

[Install]
WantedBy=multi-user.target
