# packages/greenbone-feed-sync/greenbone-nvt-sync.service
[Unit]
Description=Synchronize the greenbon nvt entries
Requires=network-online.target
After=network-online.target
After=redis.service

[Service]
Type=simple
User=gvm
Group=gvm
ExecStart=/usr/bin/sh -c '/usr/bin/greenbone-nvt-sync && sudo openvas --update-vt-info'

[Install]
WantedBy=multi-user.target
