# packages/greenbone-feed-sync/greenbone-scapdata-sync.service
[Unit]
Description=Synchronize the greenbon scapdata
Requires=network-online.target
After=network-online.target

[Service]
Type=simple
User=gvm
Group=gvm
ExecStart=/usr/bin/greenbone-scapdata-sync --rsync

[Install]
WantedBy=multi-user.target
