# packages/greenbone-feed-sync/greenbone-feed-sync.timer
[Unit]
Description=Synchronize the greenbon feed daily

[Timer]
OnCalendar=daily

[Install]
WantedBy=multi-user.target
