# packages/greenbone-feed-sync/greenbone-feed-sync.service
[Unit]
Description=Synchronize the greenbon feed
Requires=network-online.target
After=network-online.target

[Service]
Type=simple
User=gvm
Group=gvm
ExecStart=/usr/bin/greenbone-feed-sync --type GVMD_DATA

[Install]
WantedBy=multi-user.target
