# packages/greenbone-feed-sync/greenbone-scapdata-sync.timer
[Unit]
Description=Synchronize the greenbon scapdata daily

[Timer]
OnCalendar=daily

[Install]
WantedBy=multi-user.target
