# packages/3proxy/3proxy.service
[Unit]
Description=3proxy
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/3proxy /etc/3proxy/3proxy.cfg
User=threeproxy
Group=threeproxy
Restart=on-failure
PrivateUsers=true

[Install]
WantedBy=multi-user.target

