# packages/artillery/artillery.service
[Unit]
Description=Artillery Service
After=network.target

[Service]
User=nobody
ExecStart=/usr/bin/python /usr/share/artillery/artillery.py
Restart=on-abort

[Install]
WantedBy=multi-user.target
