# packages/blackarch-scanner/openvas/openvas.service
[Unit]
Description=Open Vulnerability Assessment System Scanner Daemon
Documentation=http://www.openvas.org/
After=redis.service
Requires=redis.service

[Service]
ExecStart=/usr/bin/openvassd -f

[Install]
WantedBy=multi-user.target
