# packages/gsa/gsad.service
[Unit]
Description=Greenbone Security Assistant (OpenVAS Web Frontend)
Documentation=http://www.openvas.org/
Requires=openvas-manager.service
After=openvas-manager.service

[Service]
Type=simple
ExecStart=/usr/bin/gsad -f --http-only --listen=127.0.0.1 --port=9392

[Install]
WantedBy=multi-user.target
