# packages/gsa/gsad.service
[Unit]
Description=Greenbone Security Assistent
After=ospd-openvas.service
Requires=ospd-openvas.service

[Service]
Type=simple
User=gvm
Group=gvm
EnvironmentFile=-/etc/default/gsad
ExecStart=/usr/bin/sudo /usr/bin/gsad $OPTIONS

[Install]
WantedBy=multi-user.target
