# packages/distccd-blackarch/distccd-armv7h.service
[Unit]
Description=A distributed C/C++ compiler for armv7h
Documentation=man:distccd(1)

[Service]
EnvironmentFile=/etc/conf.d/distccd-armv7h
ExecStart=/usr/bin/distccd-armv7h --daemon --no-detach $DISTCC_ARGS

[Install]
WantedBy=multi-user.target
