# packages/distccd-blackarch/distccd-armv6h.service
[Unit]
Description=A distributed C/C++ compiler for armv6h
Documentation=man:distccd(1)

[Service]
EnvironmentFile=/etc/conf.d/distccd-armv6h
ExecStart=/usr/bin/distccd-armv6h --daemon --no-detach $DISTCC_ARGS

[Install]
WantedBy=multi-user.target
